# services/car-hud-music.service
[Unit]
Description=Car-HUD Bluetooth Music Service
After=bluetooth.target

[Service]
ExecStart=/usr/bin/python3 /home/chrismslist/northstar/music_service.py
Restart=always
RestartSec=10
User=chrismslist
StandardOutput=inherit
StandardError=inherit

[Install]
WantedBy=multi-user.target
